# staging_dir/target-arm_cortex-a7_uClibc-1.0.14_eabi/root-ipq806x/etc/hotplug.d/iface/17-gl-dnsmasq
#!/bin/sh

. /lib/functions.sh
. /lib/functions/network.sh

if [  "$INTERFACE" = "guest" ];then
        if [ "$ACTION" == "ifup" ];then
		/etc/init.d/dnsmasq restart
        fi
fi
