# staging_dir/target-arm_cortex-a7_uClibc-1.0.14_eabi/root-ipq806x/etc/hotplug.d/iface/17-vpn-user
#!/bin/sh

INTERFACES="wan wan1 wwan modem tethering wg ovpn"
tmpinterface=`echo $INTERFACE | cut -f1 -d\_`
[ -f "/etc/vpn.user" ] && {
	for interface in ${INTERFACES}; do
		[ "$interface" = "$INTERFACE" -o "$interface" = "$tmpinterface" ] && {
			env -i ACTION="$ACTION" INTERFACE="$INTERFACE" DEVICE="$DEVICE" \
				/bin/sh /etc/vpn.user
		}
	done
}

exit 0
