# staging_dir/target-arm_cortex-a7_uClibc-1.0.14_eabi/root-ipq806x/etc/uci-defaults/101_open_offload
#!/bin/sh

uci set firewall.@defaults[0].flow_offloading=1
uci set firewall.@defaults[0].flow_offloading_hw=1
uci commit firewall
