# staging_dir/target-arm_cortex-a7_uClibc-1.0.14_eabi/root-ipq806x/etc/uci-defaults/99-disable-upnpd
[ -f /etc/config/upnpd ] || exit 0

uci set upnpd.config.enable_upnp=0
uci set upnpd.config.enable_natpmp=0
uci commit upnpd

[ -f /etc/init.d/miniupnpd ] && /etc/init.d/miniupnpd disable || rm /etc/config/upnpd
