# staging_dir/target-arm_cortex-a7_uClibc-1.0.14_eabi/root-ipq806x/etc/uci-defaults/gl-s2s.defaults
#!/bin/sh
uci -q batch <<-EOF >/dev/null
	delete firewall.gls2s
	set firewall.gls2s=include
	set firewall.gls2s.type=script
	set firewall.gls2s.path=/var/etc/gls2s.include
	set firewall.gls2s.reload=1
	commit firewall
EOF
exit 0
