# staging_dir/target-arm_cortex-a7_uClibc-1.0.14_eabi/root-ipq806x/etc/init.d/clear_flag
#!/bin/sh /etc/rc.common

START=99

clear_first_time_flag() {
	rm /tmp/first_time
}

start() {
	clear_first_time_flag
	
	/etc/init.d/clear_flag disable
}

