# staging_dir/target-arm_cortex-a7_uClibc-1.0.14_eabi/root-ipq806x/etc/init.d/gl_dnscrypt
#!/bin/sh /etc/rc.common

START=99

start(){
    dnscrypt=`uci get glconfig.general.dnscrypt_proxy 2>/dev/null`
    [ "$dnscrypt" == "1" ] && /etc/init.d/dnscrypt-proxy start
}
