# staging_dir/target-arm_cortex-a7_uClibc-1.0.14_eabi/root-ipq806x/etc/init.d/gl_mv1000
#!/bin/sh /etc/rc.common
. /lib/functions/gl_util.sh

START=10

start(){
	model=$(get_model)

	[ "$model" = "mv1000" ] && set_lan_mac mv1000
}
