# staging_dir/target-arm_cortex-a7_uClibc-1.0.14_eabi/root-ipq806x/etc/init.d/gl_ssids
#!/bin/sh /etc/rc.common

. /lib/functions.sh

START=30
USE_PROCD=1

ssids_cb() {
        uci set ssids.$1.disabled='0'
        uci set ssids.$1.trycount='0'
}

config_load ssids
config_foreach ssids_cb network
uci commit ssids
rm /etc/config/ssids_blacklist -rf
