# staging_dir/target-arm_cortex-a7_uClibc-1.0.14_eabi/root-ipq806x/etc/hotplug.d/ieee80211/10-wifi-detect
#!/bin/sh

[ "${ACTION}" = "add" ] && {

	[ ! -f "/etc/config/wireless" ] && exit 0

	[ -f "/tmp/gl_mv1000_wifi.pid" ] && exit 0

	/etc/init.d/gl_mv1000_wifi start "add" 
}
